# src/etc/rc.d/dhcp6c
#!/bin/sh

daemon="/usr/local/sbin/dhcp6c"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_check() {
	kill -0 $(cat /var/run/dhcp6c.pid)
}

rc_stop() {
	kill -TERM $(
		cat /var/run/dhcp6c.pid
		rm -f /var/run/dhcp6c.pid 2> /dev/null > /dev/null
	)
}

rc_cmd $1
